# node/supervisord.conf
[program:original_entrypoint]
command=npm start 
autostart=true
autorestart=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
